# r2x1-basic/files/etc/uci-defaults/99_setup_dhcp
#!/bin/sh
#
# Copyright (c) 2023 LENAR
#

uci del dhcp.lan.ra_slaac
uci commit dhcp

exit 0
