# r2x1-basic/files/etc/uci-defaults/99_setup_ncm
#!/bin/sh
#
# Copyright (c) 2023 LENAR
#

sed -i '/gcom/s/manufacturer=.*/manufacturer="notion"/g' /lib/netifd/proto/ncm.sh

exit 0
