# r2x1-basic/files/etc/uci-defaults/99_setup_dhcp
#!/bin/sh
#
# Copyright (c) 2023 LENAR
#

uci -q batch <<-EOF >/dev/null
  del dhcp.lan.ra
  del dhcp.lan.ra_flags
  del dhcp.lan.ra_slaac
  del dhcp.lan.dhcpv6
EOF
uci commit dhcp

exit 0
